# templates/default/supervisor/security_monkey.ini.erb
[program:securitymonkey]
user=<%= node['security_monkey']['user'] %>
command=python <%= node['security_monkey']['basedir'] %>/manage.py run_api_server
environment=PYTHON_EGG_CACHE="<%= node['security_monkey']['basedir'] %>/.python-eggs",SECURITY_MONKEY_SETTINGS="<%= node['security_monkey']['basedir'] %>/env-config/config-deploy.py"
autostart=true
autorestart=true
logfile=/var/log/securitymonkey/securitymonkey.log
loglevel=trace

[program:securitymonkeyscheduler]
user=<%= node['security_monkey']['user'] %>
command=python <%= node['security_monkey']['basedir'] %>/manage.py start_scheduler
directory=<%= node['security_monkey']['basedir'] %>/
environment=PYTHON_EGG_CACHE="<%= node['security_monkey']['basedir'] %>/.python-eggs",PYTHONPATH='<%= node['security_monkey']['basedir'] %>/',SECURITY_MONKEY_SETTINGS="<%= node['security_monkey']['basedir'] %>/env-config/config-deploy.py"
user=<%= node['security_monkey']['user'] %>
autostart=true
autorestart=true
logfile=/var/log/securitymonkey/securitymonkey-scheduler.log
loglevel=trace
